# nclip.service
[Unit]
Description=NClip Daemon - Clipboard History Manager
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/nclipdaemon
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0

[Install]
WantedBy=default.target